# api_service/crontab
# Run the main.py script
*/2 * * * * /usr/local/bin/python /usr/src/api_service/main.py >> /proc/1/fd/1 2>> /proc/1/fd/2
